# as04/81-ttshivhu-usb.rules
SUBSYSTEM=="input", SUBSYSTEMS=="usb", ACTION=="add", RUN+="/sbin/insmod /root/little/as04/main.ko"
